# com.pacefinder.listener.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pacefinder.listener</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>__LISTENER_PATH__</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/__USER__/Library/Logs/pacefinder.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/__USER__/Library/Logs/pacefinder.log</string>

    <key>WorkingDirectory</key>
    <string>__WORKING_DIR__</string>
</dict>
</plist>
