# simtelemetry.service
[Unit]
Description=SimTelemetry UDP Listener
After=network.target local-fs.target
# Listener falls back to local storage if USB isn't mounted, so no hard USB dep.

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/simtelemetry
ExecStart=/usr/bin/python3 /home/pi/simtelemetry/listener.py
Restart=always
RestartSec=5
# Never stop restarting regardless of how many failures occur.
StartLimitIntervalSec=0
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
